# ansible/roles/install/files/supervisor_nginx.conf
[program:nginx]
command=nginx
stdout_logfile=/var/log/supervisor/supervisor_%(program_name)s.log
stderr_logfile=/var/log/supervisor/supervisor_%(program_name)s.log
autostart=true
autorestart=true

